# build/homecontrol.service
[Unit]
Description=HomeControl Service
DefaultDependencies=no
Wants=network.target
After=network.target

[Service]
Type=simple
Restart=always

# Restart service after 10 seconds if dotnet service crashes
RestartSec=10

SyslogIdentifier=homecontrol
ExecStart=/var/homecontrol/hc-rs
WorkingDirectory=/var/homecontrol

[Install]
WantedBy = multi-user.target
